# agent/systemd/srm-monitor-agent.timer
[Unit]
Description=SRM Monitoring Agent Timer

[Timer]
OnBootSec=30
OnUnitActiveSec=1min
RandomizedDelaySec=10

[Install]
WantedBy=timers.target
