# apps/agent/systemd/adminhelper-agent.service
[Unit]
Description=AdminHelper Agent — FRPC Sync + Monitoring
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/adminhelper-agent run --once
TimeoutStartSec=60
