# apps/agent/deb/DEBIAN/postinst
#!/bin/sh
set -e
systemctl daemon-reload
mkdir -p /etc/frp/pki
chmod 700 /etc/frp/pki
mkdir -p /etc/adminhelper
chmod 700 /etc/adminhelper

# Alte systemd-Units aufraeumen (falls von frueheren Paketen vorhanden)
for unit in srm-frpc-sync.timer srm-frpc-sync.service srm-monitor-agent.timer srm-monitor-agent.service srm-agent.timer srm-agent.service; do
    if [ -f "/etc/systemd/system/${unit}" ]; then
        systemctl stop "${unit}" 2>/dev/null || true
        systemctl disable "${unit}" 2>/dev/null || true
        rm -f "/etc/systemd/system/${unit}"
    fi
done
systemctl daemon-reload

# Timer aktivieren: laeuft ohne Config harmlos (Sync/Push skippen still).
# frpc.service nur reaktivieren, wenn bereits provisioniert (frpc.toml existiert).
# Schuetzt vor Restart-Loop auf frisch installierten Systemen ohne frpc init.
if [ "$1" = "configure" ]; then
    systemctl enable --now adminhelper-agent.timer >/dev/null 2>&1 || true
    if [ -f /etc/frp/frpc.toml ]; then
        systemctl enable --now frpc.service >/dev/null 2>&1 || true
    fi
fi

echo "adminhelper-agent installiert. Timer ist aktiv (laeuft alle 5 Minuten)."
echo "FRPC Einrichtung:    sudo adminhelper-agent frpc init --url <ADMINHELPER_URL> --token <TOKEN> --server-id <ID>"
echo "Monitor Einrichtung: sudo adminhelper-agent monitor init --url <MONITOR_URL> --api-key <KEY> --server-id <ID>"
