# apps/agent/deb/DEBIAN/prerm
#!/bin/sh
set -e

# Bei Upgrade: nur stoppen, enabled-Status bleibt erhalten
# Bei Remove/Deconfigure: stop + disable
case "$1" in
    remove|deconfigure)
        systemctl stop adminhelper-agent.timer 2>/dev/null || true
        systemctl stop frpc.service 2>/dev/null || true
        systemctl disable adminhelper-agent.timer 2>/dev/null || true
        systemctl disable frpc.service 2>/dev/null || true
        ;;
    upgrade|failed-upgrade)
        systemctl stop adminhelper-agent.timer 2>/dev/null || true
        systemctl stop frpc.service 2>/dev/null || true
        ;;
esac

# Alte Units aus Vorgaenger-Paketen komplett entfernen (immer)
systemctl stop srm-agent.timer 2>/dev/null || true
systemctl stop srm-frpc-sync.timer 2>/dev/null || true
systemctl stop srm-monitor-agent.timer 2>/dev/null || true
systemctl disable srm-agent.timer 2>/dev/null || true
systemctl disable srm-frpc-sync.timer 2>/dev/null || true
systemctl disable srm-monitor-agent.timer 2>/dev/null || true
