# agent-go/systemd/frpc.service
[Unit]
Description=FRP Client (managed by SRM)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/frpc -c /etc/frp/frpc.toml
Restart=on-failure
RestartSec=5
LimitNOFILE=1048576

[Install]
WantedBy=multi-user.target
