# agent-go/systemd/srm-agent.service
[Unit]
Description=SRM Agent — FRPC Sync + Monitoring
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/srm-agent run --once
TimeoutStartSec=60
